# fsmods/etc/update-motd.d/90-footer
#!/bin/sh
#
#    90-footer - create the footer of the MOTD
#

echo
echo "Have a nive day."
echo

